# modules/tugg/templates/tugg-supervisor.cfg
[program: tugg]
directory = /opt/tugg/gigs/
username = tugg
command = /opt/tugg/gigs/bin/django run_gunicorn
stdout_logfile = /var/log/supervisor/logfile.log
stderr_logfile = /var/log/supervisor/error-logfile.log
